# server/shoplist.service
[Unit]
Description=ShopList backend (mobrauntech)
After=network.target

[Service]
Type=simple
# Adjust User and WorkingDirectory to where you clone this on the VM
User=ubuntu
WorkingDirectory=/home/ubuntu/shoplist/server
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=3
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
